# com.parceria.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.parceria.monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/miyukikasahara/documents/Parceria Dashboard/.venv/bin/python3</string>
    <string>/Users/miyukikasahara/documents/Parceria Dashboard/monitor.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/miyukikasahara/documents/Parceria Dashboard</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/miyukikasahara/documents/Parceria Dashboard/logs/launchd_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/miyukikasahara/documents/Parceria Dashboard/logs/launchd_stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
